# vendor/github.com/dgryski/carbonzipper/contrib/deb/carbonzipper.conf
description     "carbonserver proxy for graphite-web and carbonapi"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/carbonzipper -config /etc/carbonzipper/carbonzipper.conf
